# mqtt-sys-mon.service.example
[Unit]
Description=MQTT System Monitor for Home Assistant
After=network.target

[Service]
Type=simple
User=root
# Update WorkingDirectory to your installation directory
WorkingDirectory=/path/to/mqtt-sys-mon

# --- EXECUTION OPTIONS ---

# Option 1: If using system-wide python package dependencies
ExecStart=/usr/bin/python3 /path/to/mqtt-sys-mon/mqtt_sys_mon.py

# Option 2: If using a virtual environment (.venv)
# ExecStart=/path/to/mqtt-sys-mon/venv/bin/python3 /path/to/mqtt-sys-mon/mqtt_sys_mon.py

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
